# restic-backup.timer
[Unit]
Description=Restic Backup alle 1 Minute

[Timer]
OnUnitActiveSec=60
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
